# certhoover.service
[Unit]
Description=certhoover certificate transparency to ClickHouse service
After=network.target

[Service]
Type=simple
# Set the working directory to your project folder
WorkingDirectory=/var/lib/certhoover/
# Use systemd runtime directory for liveness file
RuntimeDirectory=certhoover
# Adjust the path to your binary if needed
ExecStart=/var/lib/certhoover/certhoover
# Environment variable for config override (optional)
# Environment=CERTHOOVER_CONFIG=/etc/certhoover/config.toml
Restart=on-failure
# Liveness file will be at /run/certhoover/liveness

# Liveness monitoring using systemd's Watchdog
WatchdogSec=60
WatchdogDevice=liveness.txt

[Install]
WantedBy=multi-user.target
